# scripts/com.marcelo.secretary-triage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Daily self-learning run. Install with:
      cp scripts/com.marcelo.secretary-triage.plist ~/Library/LaunchAgents/
      launchctl unload ~/Library/LaunchAgents/com.marcelo.secretary-triage.plist 2>/dev/null
      launchctl load  ~/Library/LaunchAgents/com.marcelo.secretary-triage.plist
  Run it once, right now, without waiting for 09:00:
      launchctl start com.marcelo.secretary-triage
  Stop it for good:
      launchctl unload ~/Library/LaunchAgents/com.marcelo.secretary-triage.plist
  Watch it:
      tail -f ~/Library/Logs/secretary-triage.log
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marcelo.secretary-triage</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/marceloabritta/Library/CloudStorage/GoogleDrive-marceloabritta@gmail.com/My Drive/Claude Projects/Personal Whatsapp AI/scripts/self-learning-daily.sh</string>
    </array>

    <!-- 09:00 every day. If the Mac is asleep at 09:00, launchd runs it on the next wake. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Do NOT run on load/login: installing the job should not fire a triage run. -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/marceloabritta/Library/Logs/secretary-triage.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/marceloabritta/Library/Logs/secretary-triage.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/marceloabritta/Library/CloudStorage/GoogleDrive-marceloabritta@gmail.com/My Drive/Claude Projects/Personal Whatsapp AI</string>
</dict>
</plist>
